# athena.service
[Unit]
Description=AthenaOSINT Web Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/AthenaOSINT
Environment="PATH=/root/AthenaOSINT/venv/bin"
ExecStart=/root/AthenaOSINT/venv/bin/gunicorn -k eventlet -w 1 -b 127.0.0.1:5000 --timeout 300 --max-requests 1000 --max-requests-jitter 50 run_web:socketio
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
